# roles/versitygw/templates/versitygw.service.j2
# Ansible managed: Do NOT edit this file manually!
[Unit]
Description=VersityGW S3-compatible gateway
After=network.target

[Service]
Type=simple
User={{ versitygw_system_user }}
Group={{ versitygw_system_group }}

EnvironmentFile={{ versitygw_env_file }}

ExecStart=/usr/bin/versitygw \
  --port {{ versitygw_s3_listen_addr }} \
  --webui {{ versitygw_webui_listen_addr }} \
{% if versitygw_webui_gateway_url | length > 0 %}
  --webui-gateways {{ versitygw_webui_gateway_url }} \
{% endif %}
  --iam-dir {{ versitygw_iam_dir }} \
  --health {{ versitygw_health_path }} \
  posix {{ versitygw_data_dir }}

SyslogIdentifier=versitygw
Restart=always
RestartSec=5

PrivateTmp=yes
ProtectHome=yes
NoNewPrivileges=yes
ProtectSystem=full
ProtectControlGroups=true
ProtectKernelModules=true
ProtectKernelTunables=yes

[Install]
WantedBy=multi-user.target
