# roles/victoriametrics/templates/victoriametrics.service.j2
# Ansible managed: Do NOT edit this file manually!

[Unit]
Description=VictoriaMetrics time-series database
After=network.target

[Service]
Type=simple
User={{ victoriametrics_system_user }}
Group={{ victoriametrics_system_group }}
LimitNOFILE={{ victoriametrics_max_open_files }}

ExecStart=/usr/local/bin/victoriametrics \
{% for flag, value in victoriametrics_service_args.items() %}
  -{{ flag }}={{ value }}{% if not loop.last %} \{% endif %}

{% endfor %}

SyslogIdentifier=victoriametrics
Restart=always
RestartSec=5

PrivateTmp=yes
ProtectHome=yes
NoNewPrivileges=yes
ProtectSystem=full
ProtectControlGroups=true
ProtectKernelModules=true
ProtectKernelTunables=yes

[Install]
WantedBy=multi-user.target
