# deploy/files/n8n-stack.service
[Unit]
Description=n8n Workflow Automation Stack
After=network-online.target
Wants=network-online.target

[Service]
Type=oneshot
RemainAfterExit=yes
WorkingDirectory=/root/n8n-n150
ExecStart=/usr/bin/podman-compose -f n8n-compose.yml up -d
ExecStop=/usr/bin/podman-compose -f n8n-compose.yml down
TimeoutStartSec=300

[Install]
WantedBy=multi-user.target
