# systemd/n8n-backup.service
[Unit]
Description=n8n Backup Service
After=network-online.target
Wants=network-online.target

[Service]
Type=oneshot
User=root
Group=root
WorkingDirectory=/opt/n8n
ExecStart=/opt/n8n/scripts/backup-n8n.sh
StandardOutput=journal
StandardError=journal

[Install]
WantedBy=multi-user.target
