# roles/restic/templates/restic.timer.j2
# {{ ansible_managed }}
[Unit]
Description=Restic backup timer: {{ item.key }}

[Timer]
OnCalendar={{ item.value.schedule }}
Persistent=true

[Install]
WantedBy=timers.target
