# roles/restic/templates/systemd/check.timer.j2
# {{ ansible_managed }}
[Unit]
Description=Restic check timer: {{ item.key }}

[Timer]
OnCalendar={{ item.value.check.schedule }}
Persistent=true

[Install]
WantedBy=timers.target
