# roles/victorialogs/templates/victorialogs.service.j2
# Ansible managed: Do NOT edit this file manually!

[Unit]
Description=VictoriaLogs log database ({{ victorialogs_instance_name }})
After=network.target

[Service]
Type=simple
User={{ victorialogs_system_user }}
Group={{ victorialogs_system_group }}
LimitNOFILE={{ victorialogs_max_open_files }}

ExecStart=/usr/local/bin/victoria-logs-prod \
  -httpListenAddr=127.0.0.1:{{ victorialogs_http_listen_port }} \
  -storageDataPath={{ victorialogs_data_dir }} \
  -retentionPeriod={{ victorialogs_retention_period }}

SyslogIdentifier={{ victorialogs_instance_name }}
Restart=always
RestartSec=5

PrivateTmp=yes
ProtectHome=yes
NoNewPrivileges=yes
ProtectSystem=full
ProtectControlGroups=true
ProtectKernelModules=true
ProtectKernelTunables=yes

[Install]
WantedBy=multi-user.target
